# skills/docker-local-dev/assets/templates/supervisor/python-cron.conf
; Python Cron/Scheduled Tasks via Supervisor Configuration Template
; Generated by docker-local-dev skill
; For Python applications without Celery (using APScheduler or cron)

; Option 1: APScheduler (recommended for in-process scheduling)
; pip install apscheduler

[program:python-scheduler]
command=python /app/scripts/scheduler.py
directory=/app
user=www-data
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/python-scheduler-stdout.log
stderr_logfile=/var/log/supervisor/python-scheduler-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3
environment=
    PYTHONUNBUFFERED="1",
    DJANGO_SETTINGS_MODULE="%(ENV_DJANGO_SETTINGS_MODULE)s"

; Option 2: Django management command for cron tasks
; [program:django-cron]
; command=python manage.py runcrons
; directory=/app
; user=www-data
; autostart=true
; autorestart=true
; startsecs=5
; stdout_logfile=/var/log/supervisor/django-cron-stdout.log
; stderr_logfile=/var/log/supervisor/django-cron-stderr.log
; environment=PYTHONUNBUFFERED="1",DJANGO_SETTINGS_MODULE="%(ENV_DJANGO_SETTINGS_MODULE)s"

; Option 3: Simple loop-based scheduler
; [program:simple-scheduler]
; command=python /app/scripts/simple_scheduler.py
; directory=/app
; user=www-data
; autostart=true
; autorestart=true
; startsecs=5
; stdout_logfile=/var/log/supervisor/simple-scheduler-stdout.log
; stderr_logfile=/var/log/supervisor/simple-scheduler-stderr.log
; environment=PYTHONUNBUFFERED="1"

; Example APScheduler script (scheduler.py):
; from apscheduler.schedulers.blocking import BlockingScheduler
; from apscheduler.triggers.cron import CronTrigger
;
; scheduler = BlockingScheduler()
;
; @scheduler.scheduled_job('interval', minutes=1)
; def job_every_minute():
;     print('Job running every minute')
;
; @scheduler.scheduled_job(CronTrigger(hour=0, minute=0))
; def daily_job():
;     print('Daily job at midnight')
;
; @scheduler.scheduled_job(CronTrigger(hour='*/2'))
; def job_every_2_hours():
;     print('Job running every 2 hours')
;
; if __name__ == '__main__':
;     print('Starting scheduler...')
;     scheduler.start()

; Example simple_scheduler.py (without external dependencies):
; import time
; import schedule
;
; def job():
;     print("Running scheduled task...")
;
; schedule.every(10).minutes.do(job)
; schedule.every().hour.do(job)
; schedule.every().day.at("10:30").do(job)
;
; while True:
;     schedule.run_pending()
;     time.sleep(1)
