# plugins/devops-skills/skills/docker-local-dev/assets/templates/supervisor/node-pm2.conf
; PM2 via Supervisor Configuration Template
; Generated by docker-local-dev skill
; For Node.js applications using PM2 as process manager

[program:pm2]
command=pm2-runtime start ecosystem.config.js
directory=/app
user=node
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/pm2-stdout.log
stderr_logfile=/var/log/supervisor/pm2-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile_backups=5
environment=NODE_ENV="%(ENV_NODE_ENV)s",HOME="/home/node"

; Alternative: Direct node process without PM2
; [program:node-app]
; command=node /app/dist/main.js
; directory=/app
; user=node
; autostart=true
; autorestart=true
; startsecs=10
; startretries=3
; stdout_logfile=/var/log/supervisor/node-stdout.log
; stderr_logfile=/var/log/supervisor/node-stderr.log
; environment=NODE_ENV="production"

; PM2 ecosystem.config.js example:
; module.exports = {
;   apps: [{
;     name: 'app',
;     script: 'dist/main.js',
;     instances: 'max',
;     exec_mode: 'cluster',
;     env: {
;       NODE_ENV: 'production'
;     }
;   }]
; };
