# plugins/devops-skills/skills/docker-local-dev/assets/templates/supervisor/node-cron.conf
; Node.js Cron Jobs via Supervisor Configuration Template
; Generated by docker-local-dev skill
; For running scheduled tasks in Node.js applications

; Option 1: Using node-cron library (recommended for simple tasks)
; Install: npm install node-cron
; Create a cron.js file that uses node-cron internally

[program:node-cron]
command=node /app/scripts/cron.js
directory=/app
user=node
autostart=true
autorestart=true
startsecs=5
startretries=3
stdout_logfile=/var/log/supervisor/node-cron-stdout.log
stderr_logfile=/var/log/supervisor/node-cron-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3
environment=NODE_ENV="%(ENV_NODE_ENV)s"

; Option 2: Using Agenda.js for job scheduling (recommended for complex tasks)
; [program:agenda]
; command=node /app/scripts/agenda-worker.js
; directory=/app
; user=node
; autostart=true
; autorestart=true
; startsecs=5
; startretries=3
; stdout_logfile=/var/log/supervisor/agenda-stdout.log
; stderr_logfile=/var/log/supervisor/agenda-stderr.log
; environment=NODE_ENV="%(ENV_NODE_ENV)s"

; Option 3: Using Bull Queue with scheduled jobs
; [program:bull-scheduler]
; command=node /app/scripts/scheduler.js
; directory=/app
; user=node
; autostart=true
; autorestart=true
; startsecs=5
; startretries=3
; stdout_logfile=/var/log/supervisor/bull-scheduler-stdout.log
; stderr_logfile=/var/log/supervisor/bull-scheduler-stderr.log
; environment=NODE_ENV="%(ENV_NODE_ENV)s"

; Example cron.js using node-cron:
; const cron = require('node-cron');
;
; // Run every minute
; cron.schedule('* * * * *', () => {
;   console.log('Running task every minute');
; });
;
; // Run every hour
; cron.schedule('0 * * * *', () => {
;   console.log('Running hourly task');
; });
;
; // Run every day at midnight
; cron.schedule('0 0 * * *', () => {
;   console.log('Running daily task');
; });
;
; console.log('Cron jobs started');
