# plugins/devops-skills/skills/docker-local-dev/assets/templates/supervisor/python-celery.conf
; Python Celery Workers via Supervisor Configuration Template
; Generated by docker-local-dev skill
; For Django, FastAPI, and Flask applications using Celery

[program:celery-worker]
command=celery -A %(ENV_PROJECT_NAME)s worker -l info --concurrency=2
directory=/app
user=www-data
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=600
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/celery-worker-stdout.log
stderr_logfile=/var/log/supervisor/celery-worker-stderr.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile_backups=5
environment=
    PYTHONUNBUFFERED="1",
    DJANGO_SETTINGS_MODULE="%(ENV_DJANGO_SETTINGS_MODULE)s",
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    REDIS_URL="%(ENV_REDIS_URL)s"

; Celery Beat (Scheduler) - for periodic tasks
[program:celery-beat]
command=celery -A %(ENV_PROJECT_NAME)s beat -l info --scheduler django_celery_beat.schedulers:DatabaseScheduler
directory=/app
user=www-data
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=30
stdout_logfile=/var/log/supervisor/celery-beat-stdout.log
stderr_logfile=/var/log/supervisor/celery-beat-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_backups=3
environment=
    PYTHONUNBUFFERED="1",
    DJANGO_SETTINGS_MODULE="%(ENV_DJANGO_SETTINGS_MODULE)s",
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    REDIS_URL="%(ENV_REDIS_URL)s"

; Flower (Celery monitoring) - Optional
; [program:celery-flower]
; command=celery -A %(ENV_PROJECT_NAME)s flower --port=5555
; directory=/app
; user=www-data
; autostart=true
; autorestart=true
; startsecs=10
; stdout_logfile=/var/log/supervisor/celery-flower-stdout.log
; stderr_logfile=/var/log/supervisor/celery-flower-stderr.log
; environment=
;     PYTHONUNBUFFERED="1",
;     DJANGO_SETTINGS_MODULE="%(ENV_DJANGO_SETTINGS_MODULE)s"

; Group all Celery processes
[group:celery]
programs=celery-worker,celery-beat
priority=999

; FastAPI with Celery example:
; command=celery -A app.celery_app worker -l info --concurrency=2
;
; Flask with Celery example:
; command=celery -A app.celery worker -l info --concurrency=2
